# agents/ninja-tap/ninja-tap.service
# zigbee-ninja passive capture agent (DESIGN.md §7.2).
# Install: see install.sh. Least-privilege — capture capability only, no host creds.
[Unit]
Description=zigbee-ninja ninja-tap passive Zigbee-coordinator capture
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=/usr/bin/python3 /opt/ninja-tap/ninja-tap.py \
  --collector ${ZN_COLLECTOR} \
  --token-file /etc/ninja-tap/token \
  --filter ${ZN_FILTER} \
  --iface ${ZN_IFACE}
EnvironmentFile=/etc/ninja-tap/agent.env
Restart=always
RestartSec=3

# Least privilege: only raw packet capture, nothing else.
AmbientCapabilities=CAP_NET_RAW CAP_NET_ADMIN
CapabilityBoundingSet=CAP_NET_RAW CAP_NET_ADMIN
NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=yes
PrivateTmp=yes
CPUQuota=25%
MemoryMax=128M

[Install]
WantedBy=multi-user.target
